# backend/config/crontab
00 18 * * * /app/mining_run.sh > /app/logs/mining_run_$(date +\%Y-\%m-\%d).log 2>&1
